# deploy/solar-dashboard.service
[Unit]
Description=Anker Solix Solar Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=solar
Group=solar
WorkingDirectory=/opt/solar-dashboard
ExecStart=/opt/solar-dashboard/venv/bin/python run.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
